# HomeGenie_Linux/Packager/DEBIAN/prerm
#!/bin/sh
sudo /etc/init.d/homegenie stop
sudo update-rc.d -f homegenie remove
sudo rm /etc/init.d/homegenie
#sudo rm -rf /usr/local/bin/homegenie
